# templates/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:glance-registry]
user=glance
command=/usr/bin/glance-registry

[program:glance-api]
user=glance
command=/usr/bin/glance-api

