# templates/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:nova-compute]
user=nova
command=/usr/bin/nova-compute --config-file=/etc/nova/nova.conf --config-file=/etc/nova/nova-compute.conf
